# contrib/packrat.service
[Unit]
Description=Packrat Backup Service
Documentation=https://github.com/logandonley/packrat
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=packrat
Group=packrat
ExecStart=/usr/local/bin/packrat daemon
Restart=always
RestartSec=10
TimeoutStopSec=90

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/tmp
PrivateTmp=true
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_NET_RAW
AmbientCapabilities=CAP_DAC_READ_SEARCH CAP_NET_RAW

# Environment
Environment=HOME=/home/packrat
Environment=XDG_CONFIG_HOME=/home/packrat/.config

[Install]
WantedBy=multi-user.target 